# linux/.config/systemd/user/wslg-fix.service
[Unit]
Description=Force full WSLg Weston compositor at login

[Service]
Type=oneshot
ExecStart=/bin/sh -c "sleep 2 && (xeyes >/dev/null 2>&1 &) && sleep 0.4 && pkill xeyes"
RemainAfterExit=yes

[Install]
WantedBy=default.target
