# linux/.config/systemd/user/wslg-autostart.service
[Unit]
Description=Start full WSLg Weston compositor automatically
After=graphical-session.target

[Service]
Type=oneshot
# This tiny command triggers the real compositor (Weston + Xwayland) instantly
ExecStart=/bin/sh -c "(xeyes >/dev/null 2>&1 &) && sleep 0.3 && pkill xeyes"
RemainAfterExit=yes

[Install]
WantedBy=default.target
